# nextcloud/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Example
# Runs example2 script
# ==============================================================================


bashio::log.info "Starting Nginx..."
if pgrep -f "[n]ginx:" > /dev/null; then
    pkill -ef [n]ginx:
    sleep 1
fi

if pgrep -f "[n]ginx:" > /dev/null; then
    pkill -9 -ef [n]ginx:
    sleep 1
fi

exec /usr/sbin/nginx -c /config/nginx/nginx.conf